# WeBookServer.service.in
[Unit]
Description=WeBookServer
Documentation=https://github.com/Light-Wizzard/QtService
After=network-online.target WeBookServer.socket

[Service]
Type=notify
NotifyAccess=exec
ExecStart=$${target.path}/$$TARGET --backend systemd
ExecReload=$${target.path}/$$TARGET --backend systemd reload
ExecStop=$${target.path}/$$TARGET --backend systemd stop
WatchdogSec=10
Restart=on-abnormal
RuntimeDirectory=$$TARGET

[Install]
#WantedBy=multi-user.target
WantedBy=default.target
